# ansible/roles/install-swarm-client/files/swarm.service
[Unit]
After=network.service

[Service]
ExecStart="/usr/local/bin/swarm.sh"

[Install]
WantedBy=default.target
